# 开机启动/navigator.service
# THIS IS A GENERATED FILE, NOT RECOMMENDED TO EDIT.

[Unit]
Description="bringup navigator"
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/navigator-start
User=vn

[Install]
WantedBy=multi-user.target
